# php-apache/8.2.28-apache/supervisor-laravel-filas.conf
[program:laravel-worker-alta]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=alta
autostart=false
autorestart=unexpected
startretries=3
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-media]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=media
autostart=false
autorestart=unexpected
startretries=3
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-baixa]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=baixa
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker_baixa.log
stopwaitsecs=60

[program:laravel-worker-unica]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-unica1]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica1
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-unica2]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica2
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-unica3]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica3
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-unica4]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica4
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-unica5]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=unica5
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-zepel]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=3 --timeout=60 --daemon --queue=zepel
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-credito]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=1 --timeout=240 --daemon --queue=credito
autostart=false
autorestart=unexpected
startretries=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-monitoramento]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=20 --backoff=900 --timeout=240 --queue=monitoramento
autostart=false
autorestart=unexpected
startretries=3
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-importacao-cadastro]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=1 --timeout=60 --daemon --queue=importacao-cadastro
autostart=false
autorestart=unexpected
startretries=3
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[program:laravel-worker-estatisticas]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/servicosmaisfrete-api/artisan queue:work --tries=1 --timeout=60 --daemon --queue=estatisticas
autostart=false
autorestart=unexpected
startretries=3
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/servicosmaisfrete-api/worker.log

[group:filas]
programs=laravel-worker-alta,laravel-worker-media,laravel-worker-baixa,laravel-worker-unica,laravel-worker-unica1,laravel-worker-unica2,laravel-worker-unica3,laravel-worker-unica4,laravel-worker-unica5,laravel-worker-zepel,laravel-worker-credito,laravel-worker-monitoramento,laravel-worker-importacao-cadastro,laravel-worker-estatisticas
priority=999
